# config/supervisord.conf
[inet_http_server]
port=*:8900        ; *:port for all iface
username=yiimp
password=supervisor
; check status from outside container with 'supervisorctl -u yiimp -p supervisor -s http://127.0.0.1:8900 status'

[supervisord]
logfile=/tmp/supervisord.log ; main log file; default $CWD/supervisord.log
logfile_maxbytes=50MB        ; max main logfile bytes b4 rotation; default 50MB
logfile_backups=10           ; # of main logfile backups; 0 means none, default 10
loglevel=info                ; log level; default info; others: debug,warn,trace
pidfile=/tmp/supervisord.pid ; supervisord pidfile; default supervisord.pid
nodaemon=true               ; start in foreground if true; default false
silent=false                 ; no logs to stdout if true; default false
minfds=1024                  ; min. avail startup file descriptors; default 1024
minprocs=200

; container internal memcached
[program:memcached]
command=/usr/bin/memcached -m 64 -p 11211 -u memcache -l 127.0.0.1 -P /var/run/memcached/memcached.pid
process_name=%(program_name)s ; process_name expr (default %(program_name)s)
numprocs=1                    ; number of processes copies to start (def 1)
directory=/tmp                ; directory to cwd to before exec (def no cwd)
umask=022                     ; umask for process (default None)
priority=999                  ; the relative start priority (default 999)
autostart=false               ; start at supervisord start (default: true)
autorestart=true              ; restart after exit
user=root                     ; run as user 1 (root)

; container internal cron
[program:cron]
command=/usr/sbin/cron -f
process_name=%(program_name)s ; process_name expr (default %(program_name)s)
numprocs=1                    ; number of processes copies to start (def 1)
directory=/tmp                ; directory to cwd to before exec (def no cwd)
umask=022                     ; umask for process (default None)
priority=999                  ; the relative start priority (default 999)
autostart=true                ; start at supervisord start (default: true)
autorestart=true              ; restart after exit
user=root                     ; run as user 1 (root)

; container internal haproxy
[program:haproxy]
command=/usr/sbin/haproxy -f /etc/yiimp/haproxy
process_name=%(program_name)s ; process_name expr (default %(program_name)s)
numprocs=1                    ; number of processes copies to start (def 1)
directory=/tmp                ; directory to cwd to before exec (def no cwd)
umask=022                     ; umask for process (default None)
priority=999                  ; the relative start priority (default 999)
autostart=false                ; start at supervisord start (default: true)
autorestart=true              ; restart after exit
user=root                     ; run as user 1 (root)

; apache web
[program:yiimp-apache]
command=apachectl -D FOREGROUND
stopsignal=6
process_name=%(program_name)s ; process_name expr (default %(program_name)s)
numprocs=1                    ; number of processes copies to start (def 1)
directory=/tmp                ; directory to cwd to before exec (def no cwd)
umask=022                     ; umask for process (default None)
priority=999                  ; the relative start priority (default 999)
autostart=true                ; start at supervisord start (default: true)
autorestart=true              ; restart after exit
user=root                     ; run as user 33 (www-data/apache)

; backend script
[program:yiimp-main]
command=bash main.sh
process_name=%(program_name)s ; process_name expr (default %(program_name)s)
numprocs=1                    ; number of processes copies to start (def 1)
directory=/var/www            ; directory to cwd to before exec (def no cwd)
umask=022                     ; umask for process (default None)
priority=999                  ; the relative start priority (default 999)
autostart=false                ; start at supervisord start (default: true)
autorestart=true              ; restart after exit
user=33                       ; run as user 33 (www-data/apache)

[program:yiimp-loop2]
command=bash loop2.sh
process_name=%(program_name)s ; process_name expr (default %(program_name)s)
numprocs=1                    ; number of processes copies to start (def 1)
directory=/var/www            ; directory to cwd to before exec (def no cwd)
umask=022                     ; umask for process (default None)
priority=999                  ; the relative start priority (default 999)
autostart=false                ; start at supervisord start (default: true)
autorestart=true              ; restart after exit
user=33                       ; run as user 33 (www-data/apache)

[program:yiimp-blocks]
command=bash blocks.sh
process_name=%(program_name)s ; process_name expr (default %(program_name)s)
numprocs=1                    ; number of processes copies to start (def 1)
directory=/var/www            ; directory to cwd to before exec (def no cwd)
umask=022                     ; umask for process (default None)
priority=999                  ; the relative start priority (default 999)
autostart=false                ; start at supervisord start (default: true)
autorestart=true              ; restart after exit
user=33                       ; run as user 33 (www-data/apache)

;stratum servers
[program:stratum-scrypt]
command=stratum /etc/yiimp/stratum/scrypt.conf
process_name=%(program_name)s ; process_name expr (default %(program_name)s)
numprocs=1                    ; number of processes copies to start (def 1)
directory=/tmp                ; directory to cwd to before exec (def no cwd)
umask=022                     ; umask for process (default None)
priority=999                  ; the relative start priority (default 999)
autostart=false               ; start at supervisord start (default: true)
autorestart=true              ; restart after exit
user=root                     ; run as user 33 (www-data/apache)
